# systemd/user/alert-battery.service
[Unit]
Description=Alert in case of low remaining battery

[Service]
Type=oneshot
ExecStart=/usr/local/bin/alert-battery

[Install]
WantedBy=graphical.target